# services/httpserver.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/httpserver.sh

[Install]
WantedBy=default.target
