# Packer/connectivity.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/bootscript.sh

[Install]
WantedBy=default.target

